# package/network/config/wifi-scripts/files/etc/hotplug.d/ieee80211/10-wifi-detect
#!/bin/sh

if [ ! -f "/etc/pcat-wlan-init-completed" ]; then
	exit 0
fi

[ "${ACTION}" = "add" ] && {\
	/sbin/wifi config
	ubus call network.wireless retry
}
